# systemd/orion-bgp-monitor.service
[Unit]
Description=ORION BGP Session Monitor
Documentation=https://github.com/luci-digital/Dell_R730_CQ5QBM2_ORION
After=network-online.target orion-gateway-monitor.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/orion-bgp-monitor
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orion-bgp-monitor

# Environment
Environment="TELUS_GATEWAY1=206.75.1.127"
Environment="TELUS_GATEWAY2=206.75.1.47"
Environment="TELUS_GATEWAY3=206.75.1.48"
Environment="LOCAL_AS=394955"
Environment="REMOTE_AS=6939"
Environment="ROUTER_VM_ID=200"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/orion

[Install]
WantedBy=multi-user.target
